# src/ga-bootstrap/etc/systemd/system/ga-bootstrap.service
[Unit]
Description=greenautarky first-boot bootstrap (install ga_manager → converge takes over)
Documentation=https://github.com/greenautarky/ga-bootstrap
After=hassos-supervisor.service network-online.target
Wants=network-online.target
Requires=hassos-supervisor.service
ConditionPathExists=!/mnt/data/.ga-bootstrapped

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=1500
# Optional override file under /mnt/data (= writable on HAOS, unlike /etc).
# Lets operators / bench tests point ga-bootstrap at non-default paths
# without rebuilding the OS image — e.g.
#   GHCR_CREDS_FILE=/mnt/data/ghcr-creds.json
#   BACKOFF_SEQUENCE=0 5 15 30
# The `-` prefix makes the file optional; missing = no-op.
EnvironmentFile=-/mnt/data/ga-bootstrap.env
ExecStart=/usr/sbin/ga-bootstrap
StandardOutput=journal+console
StandardError=journal+console

# Service is the entry-point of the whole user-facing stack — if it
# panics we want the boot to continue so an operator can SSH in and
# fix things rather than getting an unreachable device.
SuccessExitStatus=0

# v1.2.7 — auto-retry on transient failures (creds-missing, supervisor
# not ready yet, etc). Without this, operators had to manually
# `systemctl reset-failed && systemctl start` after fleet-manager
# pushed creds — the K7/K0/K17 BOSv1.2.12 cascade 2026-06-10. With
# Restart=on-failure + the existing ConditionPathExists=!marker, the
# unit retries automatically until creds appear OR StartLimitBurst is
# hit. Once the marker `/mnt/data/.ga-bootstrapped` is written, retry
# stops anyway (= Condition gate fires).
#
# RestartSec=180 (3 min) balances responsiveness vs storm-avoidance.
# StartLimitBurst=10 in 3600s (1h) caps the loop — after that operator
# must intervene so we don't hide a permanent problem behind infinite
# retries.
Restart=on-failure
RestartSec=180
StartLimitBurst=10
StartLimitIntervalSec=3600

[Install]
WantedBy=multi-user.target
